# deploy/coatheal-onboard.service
[Unit]
Description=COATHEAL Onboard Flight Software
After=network-online.target systemd-time-wait-sync.service
Wants=network-online.target
# Integration can flap a lot; don't let systemd give up.
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
# Type=notify so systemd uses our sd_notify READY=1 and WATCHDOG=1 messages.
Type=notify
NotifyAccess=main
User=coatheal
WorkingDirectory=/bexus/code/coatheal
Environment=COATHEAL_ENV=flight
EnvironmentFile=-/etc/coatheal/env
ExecStartPre=/bexus/code/coatheal/scripts/preflight_healthcheck.sh ${COATHEAL_CONFIG}
ExecStart=/bexus/code/coatheal/build/onboard/coatheal_onboard --config ${COATHEAL_CONFIG}
Restart=always
RestartSec=2
# BEXUS User Manual §5.9: hung-loop recovery. The main loop pings the watchdog
# every tick (default 1 Hz). 10 s gives 10 ticks of headroom before SIGKILL.
WatchdogSec=10
KillSignal=SIGINT
TimeoutStopSec=20
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
