# deploy/coatheal-onboard-debug.service
[Unit]
Description=COATHEAL Onboard (debug profile)
After=network-online.target systemd-time-wait-sync.service
Wants=network-online.target
# Flight and debug both bind port 5000 - only one can run at a time.
Conflicts=coatheal-onboard.service
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=notify
NotifyAccess=main
User=coatheal
WorkingDirectory=/bexus/code/coatheal
Environment=COATHEAL_ENV=debug
ExecStartPre=/bexus/code/coatheal/scripts/preflight_healthcheck.sh /bexus/code/coatheal/config/onboard.debug.ini
ExecStart=/bexus/code/coatheal/build/onboard/coatheal_onboard --config /bexus/code/coatheal/config/onboard.debug.ini
Restart=always
RestartSec=2
WatchdogSec=10
KillSignal=SIGINT
TimeoutStopSec=20
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
